# deploy/team-dashboard.service
[Unit]
Description=Torus Dashboard Platform - Team Control
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/team-dashboard
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
